# com.fabada.ai-news-digest.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fabada.ai-news-digest</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-l</string>
        <string>-c</string>
        <!-- Replace with your local filepath to run.sh -->
        <string>/Users/fabada/projects/ai-news-digest/run.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Weekday</key>
        <integer>6</integer>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/fabada/.local/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/Users/fabada/projects/ai-news-digest/logs/launchd-stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/fabada/projects/ai-news-digest/logs/launchd-stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/fabada/projects/ai-news-digest</string>
</dict>
</plist>
